# where do these go/fan_controller.service
[Unit]
Description=Dual PWM Fan Controller for Raspberry Pi
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/fan_controller.py
Restart=on-failure
StandardOutput=append:/var/log/fan_controller.log
StandardError=append:/var/log/fan_controller.log
SyslogIdentifier=fan_controller
User=pi
Group=pi

[Install]
WantedBy=multi-user.target
